# deploy/juyunbang-api.service
[Unit]
Description=JuYunBang API Service
After=network.target mariadb.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/juyunbang/backend
ExecStart=/opt/juyunbang/backend/venv/bin/gunicorn -c gunicorn_config.py run:app
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/juyunbang/backend

[Install]
WantedBy=multi-user.target
