# deploy/juyunbang-api.service
[Unit]
Description=JuYunBang API Service
After=network.target mariadb.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/juyunbang/code/backend
Environment="PATH=/opt/juyunbang/code/backend/venv/bin:/usr/bin"
Environment=PYTHONPATH=/opt/juyunbang/code/backend
ExecStart=/opt/juyunbang/code/backend/venv/bin/gunicorn -c gunicorn_config.py run:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
