# deploy/juyunbang-crawler.service
[Unit]
Description=JuYunBang Crawler Service (Playwright)
After=network.target mariadb.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/juyunbang/backend
Environment="PATH=/opt/juyunbang/backend/venv/bin:/usr/bin"
Environment=PYTHONPATH=/opt/juyunbang/backend
ExecStart=/opt/juyunbang/backend/venv/bin/python scheduler/task_scheduler.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
