# deploy/rejubang-api.service
[Unit]
Description=RejuBang API Service
After=network.target mariadb.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rejubang/backend
Environment="PATH=/opt/rejubang/backend/venv/bin:/usr/bin"
Environment=PYTHONPATH=/opt/rejubang/backend
ExecStart=/opt/rejubang/backend/venv/bin/gunicorn -c gunicorn_config.py run:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
